# .docker/cron/crontab
0 * * * * cd /app && drush @sites cron --yes | logger -t cron &

